# installation/udev/60-backlight.rules
SUBSYSTEM=="backlight", ACTION=="add",
RUN+="/bin/chmod 777 /sys/class/backlight/intel_backlight/brightness" 
